# build/supervisor_pg.conf
[program:postgres]
user            = postgres
command         = bash -c "/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/data -c config_file=/var/lib/postgresql/data/postgresql.conf -h $(python -c \"import socket;print socket.gethostbyname('$HOSTNAME')\")"
process_name    = %(program_name)s
stopsignal      = INT
autostart       = false
autorestart     = true
redirect_stderr = true
priority        = 100
stdout_logfile  = /var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile  = /var/log/supervisor/%(program_name)s/stderr.log

