# build/supervisor_pg.conf
[program:pg_data]
user            = postgres
command         = /usr/lib/postgresql/9.3/bin/postmaster -D /var/lib/postgresql/data -c "config_file=/var/lib/postgresql/data/postgresql.conf"
process_name    = %(program_name)s
stopsignal      = INT
autostart       = false
autorestart     = true
redirect_stderr = true
priority        = 100

