# build/supervisor_sshd.conf
[program:sshd]
command=bash -c "if [ !-d /var/run/sshd ];then install -d -m 0755 -o root -g root /var/run/sshd; fi; /usr/sbin/sshd -D"
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log
process_name    = %(program_name)s
autostart       = false
autorestart     = true
redirect_stderr = true

